# you-tui-daemon.service
[Unit]
Description=YouTube TUI Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/you-tui-daemon
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Run as the user who enabled the service (not needed for user services)

[Install]
WantedBy=default.target
